# systemd/infinivirt-nftables.service
[Unit]
Description=Infinivirt nftables Firewall Rules
Documentation=https://github.com/infinibay/infinivirt
After=network.target nftables.service
Before=infinibay.service
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Load persistent firewall rules on boot
ExecStart=/usr/sbin/nft -f /etc/infinivirt/nftables/infinivirt.nft

# Clear infinivirt table on stop (optional - leaves rules in place by default)
# ExecStop=/usr/sbin/nft delete table bridge infinivirt

# Retry configuration
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=no
PrivateTmp=yes
ReadOnlyPaths=/etc/infinivirt/nftables
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
